# systemd/seekarr-webui.service
[Unit]
Description=Seekarr (web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/seekarr
ExecStart=/opt/seekarr/venv/bin/python /opt/seekarr/webui_main.py --config /etc/seekarr/config.yaml --host 127.0.0.1 --port 8788
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
# Web UI can run standalone (no worker service) and will auto-run searches by default.
Environment=WEBUI_AUTORUN_DEFAULT=1
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/seekarr
EnvironmentFile=-/etc/seekarr/seekarr.env

[Install]
WantedBy=multi-user.target
