# systemd/seekarr-console.service
[Unit]
Description=Seekarr (console worker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/seekarr
ExecStart=/opt/seekarr/venv/bin/python /opt/seekarr/main.py --config /etc/seekarr/config.yaml
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/seekarr
EnvironmentFile=-/etc/seekarr/seekarr.env

[Install]
WantedBy=multi-user.target

